# crontab
*/5     *       *       *       *       rsync -avr ~/tmpfs/header_cache ~/.mutt
